# root/etc/services.d/php5.6-fpm/run
#!/usr/bin/env bash
exec /usr/sbin/php-fpm5.6 --allow-to-run-as-root -c /etc/php/5.6/fpm --nodaemonize
EOT
