# root/etc/services.d/php7.2-fpm/run
#!/usr/bin/env bash
exec /usr/sbin/php-fpm7.2 --allow-to-run-as-root -c /etc/php/7.2/fpm --nodaemonize
EOT
