# root/etc/services.d/php7.4-fpm/run
#!/usr/bin/env bash
exec /usr/sbin/php-fpm7.4 --allow-to-run-as-root -c /etc/php/7.4/fpm --nodaemonize
EOT
