# root/etc/services.d/redis-server/run
#!/usr/bin/env bash
exec /usr/bin/redis-server --daemonize no
EOT
